# service/sdboot-snaps-watch.path
[Unit]
Description=Watch for snapper snapshot changes

[Path]
PathChanged=/.snapshots
Unit=sdboot-snaps-watch.service
TriggerLimitIntervalSec=5min

[Install]
WantedBy=multi-user.target
